# DEBIAN/postinst
#!/bin/bash

PMS_LIB_DIR=/usr/lib/plexmediaserver
PMS_DOWNLOAD_DIR=/tmp/plex_tmp_download	

get_pms_build(){

	local cpuinfo=$(grep Features /proc/cpuinfo)
	if [[ $cpuinfo =~ .*neon.* ]] && [[ $cpuinfo =~ .*vfpv4.* ]] && [[ $cpuinfo =~ .*thumb.* ]] && [[ $cpuinfo =~ .*idiva.* ]]; then
		echo "netgear"		
	elif [[ $cpuinfo =~ .*neon.* ]]; then
		echo "synology_neon"
	else
		echo "synology"
	fi
}

install_netgear(){

	local PMS_URL='https://downloads.plex.tv/plex-media-server/0.9.17.2.2159-2bd156c/plexmediaserver-ros6-binaries-annapurna_0.9.17.2.2159-2bd156c_armel.deb'
	local PMS_HASH='ebb5b0a04790ed9e2035f45ebfc518ec74003fe77c535ba39e7f968a055a582d'

	echo "Downloading readynas package ..."
	cd $PMS_DOWNLOAD_DIR	
	curl --progress-bar -o readynas.deb $PMS_URL

	local PMS_DOWNLOAD_HASH=`sha256sum readynas.deb | cut -d' ' -f1`

	if [ "$PMS_HASH" != "$PMS_DOWNLOAD_HASH" ]
	then
		echo "Checksum mismatch. Downloaded file does not match this package."
		exit 1
	else
		echo "Passed checksum test."
	fi

	echo "Extracting readynas.deb ..."
        dpkg-deb --fsys-tarfile readynas.deb | tar -xf - -C $PMS_LIB_DIR/ --strip-components=4 ./apps/plexmediaserver-annapurna/Binaries

	# remove not used files
	rm $PMS_LIB_DIR/config.xml
	
}

install_synology(){

	local PMS_URL='https://downloads.plex.tv/plex-media-server/0.9.17.2.2159-2bd156c/PlexMediaServer-0.9.17.2.2159-2bd156c-arm7.spk'
	local PMS_HASH='c079b63da090660231b18226b5a1c14d1810f30b428fa987e2f26a2ed39515bd'

	if [[ $1 == "synology_neon" ]]; then
		PMS_URL='https://downloads.plex.tv/plex-media-server/0.9.17.2.2159-2bd156c/PlexMediaServer-0.9.17.2.2159-2bd156c-neon-armv7.spk'
		PMS_HASH='dc8b073f0caffd1803ebfab9a266bbb6f28db713ecbb2720d7d480898b28c63d'
	fi

	echo "Downloading synology package ..."
	cd $PMS_DOWNLOAD_DIR	
	curl --progress-bar -o synology.tar $PMS_URL

	local PMS_DOWNLOAD_HASH=`sha256sum synology.tar | cut -d' ' -f1`

	if [ "$PMS_HASH" != "$PMS_DOWNLOAD_HASH" ]
	then
		echo "Checksum mismatch. Downloaded file does not match this package."
		exit 1
	else
		echo "Passed checksum test."
	fi

	echo "Extracting synology.tar ..."
        tar -xOf synology.tar package.tgz | tar -xzf - -C $PMS_LIB_DIR/

	# remove not used files
	rm -r $PMS_LIB_DIR/dsm_config
	if [ -e $PMS_LIB_DIR/ld-linux.so.3 ]
	then
		rm $PMS_LIB_DIR/ld-linux.so.3
	fi
	
}


case "$1" in
    configure)
	adduser --quiet --system --shell /bin/bash --home /var/lib/plexmediaserver plex
	update-rc.d plexmediaserver defaults
	/etc/init.d/plexmediaserver stop

	# create dirs 
	mkdir -p $PMS_DOWNLOAD_DIR	
	mkdir -p $PMS_LIB_DIR

	pmsbuild=$(get_pms_build)

	if [[ $pmsbuild == "netgear" ]]; then
		install_netgear
	else
		install_synology $pmsbuild
	fi

	# remove tmp data
	cd /tmp
	rm -r $PMS_DOWNLOAD_DIR/

	# fix web profile
	sed -i 's/name="audio.channels" value="6"/name="audio.channels" value="2"/g' /usr/lib/plexmediaserver/Resources/Profiles/Web.xml

	/etc/init.d/plexmediaserver start
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

exit 0
