# etc/init/plexmediaserver.conf
# plexpms - service job file

description "Plex Media Server"
author "http://www.plexapp.com/"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Sets nice and ionice level for job
nice -5

# What to execute
script
    if [ -r /etc/default/plexmediaserver ]; then
        . /etc/default/plexmediaserver
    fi 
    start-stop-daemon --start -c $PLEX_MEDIA_SERVER_USER --exec /usr/sbin/start_pms
end script
